# crontab
* * * * * docker exec laravel_php sh -c " php artisan schedule:run >> /dev/null 2>&1 "
